# meta-openvuplus/recipes-vuplus/enigma2/enigma2-transtreamproxy/transtreamproxy@.service
[Unit]

[Service]
ExecStart=-/usr/bin/transtreamproxy
StandardInput=socket
